# OpenWrt/Rootfs/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

rm /root/record.time
/root/channel.sh &
python /root/sniff.py &
python /root/upload.py &
exit 0
